# crontab
0 12 * * * /usr/bin/python3 /pfad/zu/deinem/update_checker.py
